# stubs/supervisor/conf.d/echo.ini
[program:echo]
process_name = %(program_name)s_%(process_num)02d
command = /bin/bash -c "%(ENV_NPX_PATH)s %(ENV_PM_DIR)s/node_modules/.bin/laravel-echo-server start --force --dev"
numprocs = 1
user = root
autostart = true
autorestart = true
startsecs = 0
redirect_stderr = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
environment = DD_SERVICE="events",DD_ENV="processmaker-dev",DD_VERSION="4.5.1",DD_AGENT_HOST="datadog-agent"
