# stubs/supervisor/conf.d/horizon.ini
[program:horizon]
process_name = %(program_name)s_%(process_num)02d
command = /bin/bash -c "%(ENV_PHP_BINARY)s %(ENV_PM_DIR)s/artisan horizon --no-interaction"
autostart = true
autorestart = true
user = root
startsecs = 0
redirect_stderr = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
environment = DD_SERVICE="horizon",DD_ENV="processmaker-dev",DD_VERSION="4.5.1",DD_AGENT_HOST="datadog-agent",DD_TRACE_AUTO_FLUSH_ENABLED=1,DD_TRACE_GENERATE_ROOT_SPAN=0,DD_PROFILING_ENABLED=1,DD_TRACE_CLI_ENABLED=1
