# stubs/supervisor/conf.d/horizon.ini
[program:horizon]
process_name = %(program_name)s_%(process_num)02d
command = /bin/bash -c "%(ENV_PHP_BINARY)s %(ENV_PM_DIR)s/artisan horizon --no-interaction"
autostart = true
autorestart = true
user = root
startsecs = 0
stopwaitsecs=43200
redirect_stderr = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
environment = DD_PROFILING_ENABLED=1
