# stubs/supervisor/conf.d/nginx.ini
[program:nginx]
process_name = %(program_name)s_%(process_num)02d
command = /bin/bash -c "nginx -c /etc/nginx/nginx.conf -g 'daemon off;'"
numprocs = 1
user = root
autostart = true
autorestart = true
startsecs = 0
redirect_stderr = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
