# stubs/supervisor/conf.d/nginx.ini
[program:nginx]
process_name = %(program_name)s_%(process_num)02d
command = /bin/bash -c "nginx -c /etc/nginx/nginx.conf -g 'daemon off;'"
user = root
autostart = true
autorestart = true
startsecs = 0
redirect_stderr = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
environment = DD_SOURCE="nginx",DD_ENV="processmaker-dev",DD_VERSION="4.5.0",DD_AGENT_HOST="datadog-agent"
