# stubs/supervisor/conf.d/php-fpm.ini
[program:php-fpm]
process_name = %(program_name)s_%(process_num)02d
command = /bin/bash -c "%(ENV_PHP_FPM_BINARY)s --fpm-config /etc/php/%(ENV_PHP_VERSION)s/fpm/php-fpm.conf --nodaemonize --allow-to-run-as-root"
user = root
autostart = true
autorestart = true
startsecs = 0
redirect_stderr = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
environment = DD_SERVICE="web",DD_ENV="processmaker-dev",DD_VERSION="4.5.0",DD_AGENT_HOST="datadog-agent"
