# stubs/supervisor/conf.d/scheduler.ini
[program:scheduler]
process_name = %(program_name)s_%(process_num)02d
command = /bin/bash -c "%(ENV_PHP_BINARY)s %(ENV_PM_DIR)s/artisan schedule:work"
autostart = true
autorestart = true
user = root
numprocs = 1
startsecs = 0
redirect_stderr = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
environment = DD_TRACE_SPANS_LIMIT=4096
